# .config/systemd/user/redshift.service
[Unit]
Description="Redshift autostart service"
[Service]
Environment=DISPLAY=:0
ExecStart=redshift
[Install]
WantedBy=default.target
